# deploy/systemd/preston-worker.service
[Unit]
Description=Preston AI OS worker (simulation only, bounded)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Dedicated least-privilege worker identity (separate from Hermes and the owner).
User=preston-worker
Group=preston-worker
WorkingDirectory=/srv/preston-os/apps/dashboard
# Worker-only env + token store (SUPABASE_RUNTIME_TOKEN_STORE points inside it).
EnvironmentFile=/etc/preston/worker.env
# Runs the compiled dispatcher (npm run build:os-runtime) in a bounded
# simulate-loop. Execution stays disabled by controls; nothing runs live.
ExecStart=/usr/bin/node dist/os-runtime/bin.js worker-loop --max 5
Restart=no
TimeoutStartSec=120
RuntimeMaxSec=300
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Creates/owns /var/log/preston (0755, user preston) and grants write access,
# so ProtectSystem=strict does not abort startup on a fresh host.
LogsDirectory=preston
StandardOutput=append:/var/log/preston/worker.log
StandardError=append:/var/log/preston/worker.log
# No [Install]: the SERVICE never auto-starts. The timer (owner-enabled) fires it.
